# QSDM/deploy/systemd/qsdm.service.d/cors.conf
# /etc/systemd/system/QSD.service.d/cors.conf
#
# CORS allowlist for the public API surface.
#
# Browser-origin enforcement is done by pkg/api/CORSMiddleware (see
# QSD/source/pkg/api/cors.go). The middleware is deny-by-default:
# an empty allowlist means EVERY browser request to the API returns
# 403 with no `Access-Control-Allow-Origin` header — which is what
# the user sees as
#
#   Access to fetch at 'https://api.QSD.tech/...' from origin
#   'https://QSD.tech' has been blocked by CORS policy
#
# in the browser console. Without this drop-in, the static-site
# wallet (QSD.tech/wallet.html), the operator dashboard
# (dashboard.QSD.tech) and any other browser-only consumer of the
# API will fail with that exact error.
#
# Install:
#
#   sudo install -m 0644 cors.conf \
#     /etc/systemd/system/QSD.service.d/cors.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart QSD.service
#
# Verify (from the same host or from outside):
#
#   curl -sSI -X OPTIONS \
#     -H 'Origin: https://QSD.tech' \
#     -H 'Access-Control-Request-Method: GET' \
#     https://api.QSD.tech/api/v1/wallet/balance?address=DEADBEEF \
#     | grep -i '^access-control-\|^vary:'
#
# Expect:
#
#   access-control-allow-origin: https://QSD.tech
#   access-control-allow-methods: GET, POST, PUT, DELETE, PATCH, OPTIONS
#   access-control-allow-headers: Authorization, Content-Type, X-CSRF-Token, X-Timestamp, X-Nonce, X-Signature
#   access-control-max-age: 600
#   vary: Origin
#
# Bearer-token auth means the wallet does NOT need credentialed
# CORS, so we deliberately leave AllowCredentials at its `false`
# default and do NOT use the `*` wildcard (which would be ignored
# anyway under credentialed mode per the Fetch spec).
#
# Operators who add a new browser-facing surface MUST extend the
# list below and `systemctl restart QSD.service`. Disallowed
# origins continue to 403 (test surface:
# `pkg/api/cors_test.go::TestCORSRejectsDisallowedOrigin`).

[Service]
Environment="QSD_CORS_ALLOWED_ORIGINS=https://QSD.tech,https://www.QSD.tech,https://dashboard.QSD.tech"
