# QSDM/deploy/systemd/qsdm.service.d/task-registry.conf
[Service]
Environment="QSD_TASK_REGISTRY_PATH=/opt/QSD/QSD-hive-tasks.json"
Environment="QSD_TASK_ACTION_LOG_PATH=/opt/QSD/QSD_task_actions.ndjson"
