# QSDM/config/qsdm.service
[Unit]
Description=Quantum-Secure Dynamic Mesh Ledger (QSD) Node
After=network.target

[Service]
Type=simple
User=QSD
Group=QSD
WorkingDirectory=/opt/QSD
ExecStart=/opt/QSD/QSD
Restart=always
RestartSec=10

# Environment variables
Environment="CGO_ENABLED=1"
Environment="LD_LIBRARY_PATH=/opt/QSD/liboqs_install/lib:/opt/QSD/liboqs_install/lib64:/usr/local/lib64:/usr/local/lib"
Environment="CONFIG_FILE=/opt/QSD/QSD.toml"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/QSD

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=QSD

[Install]
WantedBy=multi-user.target
